# crontab
# Example run with bash shell
* * * * * root /var/www/crontab.sh


